# packaging/aegis.service
[Unit]
Description=Aegis eBPF Security Agent
After=network.target

[Service]
ExecStart=/usr/bin/aegis-agent --config /etc/aegis/config.yaml
Restart=always
User=root

[Install]
WantedBy=multi-user.target
